# rplidar/rplidar.conf
start on startup
task
script
  mkdir /var/run/rplidar
  chown turtlebot /var/run/rplidar
end script

